# config/systemd/user/allsync.service
[Unit]
Description=Allsync

[Service]
ExecStart="/home/cynerd/.local/bin/allsync"
Environment="AUTORUN=y"
