# config/systemd/user/allsync.timer
[Unit]
Description=Run allsync every two minutes

[Timer]
OnCalendar=*:0/2

[Install]
WantedBy=timers.target
