# config/systemd/user/davmail.service
[Unit]
Description=DavMail

[Service]
ExecStart="/run/current-system/sw/bin/davmail"
