# config/systemd/user/waybar.service
[Unit]
Description=Waybar

[Service]
ExecStart="/run/current-system/sw/bin/waybar"

[Install]
WantedBy="graphical-session.target"
