# config/systemd/user/wayvnc1.service
[Unit]
Description=wayvnc for wayland display 1

[Service]
ExecStart="/home/cynerd/.local/sbin/wayvnc1"
